# dev-util/hoogle/files/hoogle-generate.timer
[Unit]
Description=Run hoogle generate

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
