# dev-util/hoogle/files/hoogle.timer
[Unit]
Description=Restart hoogle periodically

[Timer]
Unit=hoogle.service
OnCalendar=daily

[Install]
WantedBy=timers.target
