# dev-util/hoogle/files/hoogle-restart.service
[Unit]
Description=restart hoogle

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart hoogle
