# dev-util/hoogle/files/hoogle-restart.timer
[Unit]
Description=Restart hoogle periodically

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
